# xdg-autostart/forge.desktop
[Desktop Entry]
Name=Forge Desktop
Comment=Forge Desktop Environment Session
Exec=forge-compositor
Type=Application
DesktopNames=Forge
